# templates/network-watchdog.service.j2
[Unit]
Description=Cluster Network Watchdog Service
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/network-watchdog.sh

[Install]
WantedBy=multi-user.target