# templates/network-watchdog.timer.j2
[Unit]
Description=Run Cluster Network Watchdog every minute

[Timer]
OnBootSec=60
OnUnitActiveSec=60

[Install]
WantedBy=timers.target